# homebotctl/homebot.service
[Unit]
Description=Homebot Service
[Service]
ExecStart={HOME}/homebot live
[Install]
WantedBy=multi-user.target
